# net/mimic/files/etc/init.d/mimic
#!/bin/sh /etc/rc.common

START=19
USE_PROCD=1

PROG=/usr/bin/mimic
KMOD_NAME=mimic

ensure_kmod_loaded() {
	[ -d "/sys/module/$KMOD_NAME" ] && return 0
	modprobe "$KMOD_NAME"
}

start_instance() {
	local cfg="$1"
	local enabled interface config_file

	config_get_bool enabled "$cfg" enabled 0
	[ "$enabled" -eq 1 ] || return 0

	config_get interface "$cfg" interface
	[ -n "$interface" ] || {
		logger -t mimic "Section '$cfg' missing 'interface', skipping"
		return 0
	}

	config_get config_file "$cfg" config_file "/etc/mimic/${interface}.conf"

	if [ ! -e "/sys/class/net/${interface}" ]; then
		logger -t mimic "Interface '$interface' does not exist, skipping (cfg=$cfg)"
		return 0
	fi

	if [ ! -r "$config_file" ]; then
		logger -t mimic "Config '$config_file' missing or unreadable, skipping (cfg=$cfg)"
		return 0
	fi

	procd_open_instance "$cfg"
	procd_set_param command "$PROG" run "$interface" -F "$config_file"
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "mimic"
}

section_match() {
	[ "$1" = "$instance" ] && instance_found=1
}

# `service mimic <action> [<name>]` — when <name> is given, only that UCI
# section is touched. rc_procd uses procd's "add" semantics in that case
# (won't disturb other instances); without a name it uses "set" (full
# replace, removed instances stop). procd_kill (default stop) honours the
# instance arg too.
start_service() {
	local instance="$1"

	if ! ensure_kmod_loaded; then
		logger -t mimic "Failed to load kernel module '$KMOD_NAME'"
		return 1
	fi

	config_load mimic

	if [ -n "$instance" ]; then
		local instance_found=0
		config_foreach section_match mimic
		if [ "$instance_found" -ne 1 ]; then
			logger -t mimic "Section '$instance' not found in /etc/config/mimic"
			return 1
		fi
		start_instance "$instance"
	else
		config_foreach start_instance mimic
	fi
}

# Mimic loads its config at startup and has no SIGHUP handler, so reload
# is implemented as stop+start. The instance arg propagates to both.
reload_service() {
	stop "$@"
	start "$@"
}
